# root/etc/s6-overlay/s6-rc.d/svc-tftp/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec /usr/sbin/in.tftpd  -L --secure  ${PORT_RANGE:+--port-range $PORT_RANGE}  /tftpboot
